# llm/q2/launchd/com.orangelabs.q2.mongo-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.orangelabs.q2.mongo-tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/tigerj/Work/OrangeLabs/Orange/Orange-AI/q2_mongo_tunnel.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/tigerj/Work/OrangeLabs/Orange/Orange-AI</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>Q2_MONGO_TUNNEL_LOCAL_PORT</key>
    <string>27018</string>
    <key>Q2_MONGO_TUNNEL_SSH_HOST</key>
    <string>golgong.com</string>
    <key>Q2_MONGO_TUNNEL_SSH_PORT</key>
    <string>10022</string>
    <key>Q2_MONGO_TUNNEL_SSH_USER</key>
    <string>orange</string>
    <key>Q2_MONGO_TUNNEL_PASSWORD</key>
    <string>Qjawlsdk12!@</string>
    <key>Q2_MONGO_TUNNEL_REMOTE_HOST</key>
    <string>127.0.0.1</string>
    <key>Q2_MONGO_TUNNEL_REMOTE_PORT</key>
    <string>27017</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/tigerj/Work/OrangeLabs/Orange/Orange-AI/llm/q2/output/q2_mongo_tunnel.launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/tigerj/Work/OrangeLabs/Orange/Orange-AI/llm/q2/output/q2_mongo_tunnel.launchd.err.log</string>
</dict>
</plist>
