# llm/q2/launchd/com.orangelabs.q2.local-lm-trainer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.orangelabs.q2.local-lm-trainer</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/tigerj/Work/OrangeLabs/Orange/Orange-AI/q2_local_lm_trainer.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/tigerj/Work/OrangeLabs/Orange/Orange-AI</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/tigerj/Work/OrangeLabs/Orange/Orange-AI/venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>Q2_PROJECT</key>
    <string>q2</string>
    <key>Q2_SERVER</key>
    <string>http://127.0.0.1:3184</string>
    <key>Q2_LOCAL_TRAINER_INTERVAL</key>
    <string>1</string>
    <key>Q2_LOCAL_TRAINER_MODE</key>
    <string>step</string>
    <key>Q2_LOCAL_TRAINER_BATCH</key>
    <string>1</string>
    <key>Q2_LOCAL_TRAINER_TIMEOUT</key>
    <string>180</string>
    <key>Q2_LOCAL_TRAINER_RETRIES</key>
    <string>3</string>
    <key>Q2_LOCAL_TRAIN_LIMIT</key>
    <string>2000</string>
    <key>Q2_LOCAL_EVAL_TRAIN_LIMIT</key>
    <string>2000</string>
    <key>Q2_LOCAL_ENABLED</key>
    <string>1</string>
    <key>Q2_OLLAMA_URL</key>
    <string>http://127.0.0.1:11434</string>
    <key>Q2_OLLAMA_MODEL</key>
    <string>qwen3-coder:30b</string>
    <key>Q2_OLLAMA_TIMEOUT</key>
    <string>180</string>
    <key>ORANGE_MONGO_HOST</key>
    <string>127.0.0.1</string>
    <key>ORANGE_MONGO_PORT</key>
    <string>27018</string>
    <key>ORANGE_MONGO_USER</key>
    <string>orange</string>
    <key>ORANGE_MONGO_PASSWORD</key>
    <string>orange_p@ss</string>
    <key>ORANGE_SOURCE_MONGO_DB</key>
    <string>public</string>
    <key>ORANGE_SOURCE_MONGO_AUTH_SOURCE</key>
    <string>public</string>
    <key>ORANGE_STORE_MONGO_DB</key>
    <string>ai</string>
    <key>ORANGE_STORE_MONGO_AUTH_SOURCE</key>
    <string>ai</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/tigerj/Work/OrangeLabs/Orange/Orange-AI/llm/q2/output/q2_local_lm_trainer.launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/tigerj/Work/OrangeLabs/Orange/Orange-AI/llm/q2/output/q2_local_lm_trainer.launchd.err.log</string>
</dict>
</plist>
